# packaging/debian/fips-health-recover.service
[Unit]
Description=Recover an unresponsive FIPS daemon
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart fips.service
