# packaging/debian/fips-healthcheck.timer
[Unit]
Description=Run the end-to-end FIPS WebRTC and FSP health probe

[Timer]
OnBootSec=30s
OnUnitInactiveSec=60s
RandomizedDelaySec=10s
AccuracySec=5s
Persistent=true
Unit=fips-healthcheck.service

[Install]
WantedBy=timers.target
